# droneserver.service
[Unit]
Description=MAVLink MCP Server - HTTP/SSE Mode for ChatGPT
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/droneserver
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Start the MCP server
ExecStart=/root/droneserver/start_http_server.sh

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=droneserver

# Security hardening (optional)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

