# ngrok.service
[Unit]
Description=ngrok HTTPS Tunnel for MAVLink MCP Server
After=network.target droneserver.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root

# Start ngrok tunnel on port 8080
# --host-header rewrites Host header to localhost:8080 (required for MCP SDK 1.23+)
ExecStart=/usr/local/bin/ngrok http 8080 --host-header="localhost:8080" --log=stdout

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ngrok

[Install]
WantedBy=multi-user.target

